# scripts/supervisord/supervisord_dispatcherd.conf
# Supervisord configuration for the dispatcherd container
# Local testing: supervisord -c scripts/supervisord/supervisord_dispatcherd.conf
# Requires: supervisord and the package installed (pip install -e .)

[supervisord]
nodaemon = true
pidfile = /tmp/supervisord_dispatcherd.pid
logfile = /dev/stdout
logfile_maxbytes = 0
childlogdir = /tmp

[program:dispatcherd]
command = metrics-service dispatcherd --workers 1 --timeout 3600 --max-tasks 100
autorestart = true
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[unix_http_server]
file = /tmp/supervisor_dispatcherd.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor_dispatcherd.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
