# scripts/supervisord/supervisord_web.conf
# Supervisord configuration for the web container (nginx + gunicorn)
# Local testing: supervisord -c scripts/supervisord/supervisord_web.conf
# Requires: nginx, supervisord, and the package installed (pip install -e .)

[supervisord]
nodaemon = true
pidfile = /tmp/supervisord_web.pid
logfile = /dev/stdout
logfile_maxbytes = 0
childlogdir = /tmp

[program:nginx]
command = nginx -c %(here)s/../nginx/nginx.conf -g "daemon off;"
autorestart = true
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[program:gunicorn]
command = sh -c 'exec python3.12 -m gunicorn metrics_service.wsgi:application --bind 127.0.0.1:8000 --workers ${GUNICORN_WORKERS:-4} --capture-output --access-logfile - --error-logfile - --log-level ${GUNICORN_LOG_LEVEL:-info}'
autorestart = true
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[unix_http_server]
file = /tmp/supervisor_web.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor_web.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
