# scripts/supervisord/supervisord_dispatcherd.conf
# Supervisord configuration for the dispatcherd container
# Local testing: supervisord -c scripts/supervisord/supervisord_dispatcherd.conf
# Requires: supervisord and the package installed (pip install -e .)

[supervisord]
nodaemon = true
pidfile = /tmp/supervisord_dispatcherd.pid
logfile = /dev/stdout
logfile_maxbytes = 0
childlogdir = /tmp

[program:dispatcherd]
command = sh -c 'exec metrics-service dispatcherd --workers ${DISPATCHERD_WORKERS:-4} --timeout ${DISPATCHERD_TIMEOUT:-3600} --max-tasks ${DISPATCHERD_MAX_TASKS:-100} --log-level ${DISPATCHERD_LOG_LEVEL:-INFO}'
autorestart = true
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[unix_http_server]
file = /tmp/supervisor_dispatcherd.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor_dispatcherd.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
