# scripts/supervisord/supervisord_scheduler.conf
# Supervisord configuration for the scheduler container
# Local testing: supervisord -c scripts/supervisord/supervisord_scheduler.conf
# Requires: supervisord and the package installed (pip install -e .)

[supervisord]
nodaemon = true
pidfile = /tmp/supervisord_scheduler.pid
logfile = /dev/stdout
logfile_maxbytes = 0
childlogdir = /tmp

[program:scheduler]
command = sh -c 'exec metrics-service scheduler --check-interval ${SCHEDULER_CHECK_INTERVAL:-60}'
autorestart = true
stopasgroup = true
killasgroup = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[unix_http_server]
file = /tmp/supervisor_scheduler.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor_scheduler.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
