# deployment/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=1001
logfile=/opt/app-root/src/logs/supervisord.log
pidfile=/opt/app-root/src/supervisord.pid
childlogdir=/opt/app-root/src/logs

[unix_http_server]
file=/opt/app-root/src/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///opt/app-root/src/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Web/API Service (uWSGI)
[program:web]
command=/opt/app-root/bin/uwsgi --ini /opt/app-root/etc/uwsgi.ini
directory=/opt/app-root/src
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/opt/app-root/src/logs/web.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=100
user=1001
environment=PATH="/opt/app-root/bin:%(ENV_PATH)s",DJANGO_SETTINGS_MODULE="metrics_service.settings"

# Task Dispatcher (dispatcherd)
[program:dispatcher]
command=/opt/app-root/bin/python manage.py run_dispatcherd --workers=4 --timeout=3600 --max-tasks=100 --log-level=INFO
directory=/opt/app-root/src
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/opt/app-root/src/logs/dispatcher.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=200
user=1001
environment=PATH="/opt/app-root/bin:%(ENV_PATH)s",DJANGO_SETTINGS_MODULE="metrics_service.settings"

# Task Scheduler (APScheduler)
[program:scheduler]
command=/opt/app-root/bin/python manage.py run_task_scheduler --log-level=INFO --check-interval=60
directory=/opt/app-root/src
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/opt/app-root/src/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
priority=300
user=1001
environment=PATH="/opt/app-root/bin:%(ENV_PATH)s",DJANGO_SETTINGS_MODULE="metrics_service.settings"
